# all-in-one/supervisor/php-fpm.conf
[program:php]
command=php-fpm -F
autostart=true
autorestart=true
stderr_logfile=/var/log/php.err.log
stdout_logfile=/var/log/php.out.log
